# crontab
* * * * * sudo a.sh
* * * * * sudo c.sh
